# dot_chezmoi_etc/systemd/system/docker.service.d/http-proxy.conf.tmpl
{{- if .proxy }}
[Service]
Environment="HTTP_PROXY={{ .proxy }}"
Environment="HTTPS_PROXY={{ .proxy }}"
Environment="NO_PROXY={{ .no_proxy }}"
{{- end }}